# .build/aria2.supervisor.conf
[supervisord]
nodaemon=false

[program:background_process]
command=aria2c --enable-rpc --save-session /cloudreve/data --listen-port 6888 --dht-listen-port 6888
autostart=true
autorestart=true